# deployment/supervisor/zkt-attendance.conf
# =============================================================================
# Supervisor Configuration for ZKTeco Attendance System
# =============================================================================
# Place this file at: /etc/supervisor/conf.d/zkt-attendance.conf
# Then run: supervisorctl reread && supervisorctl update
# =============================================================================

# Laravel Queue Worker — processes background jobs (sync, commands, etc.)
[program:zkt-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/zkt-attendance/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/zkt-attendance/storage/logs/queue-worker.log
stopwaitsecs=3600

# Laravel Scheduler — runs cron jobs every minute
[program:zkt-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do php /var/www/zkt-attendance/artisan schedule:run --no-interaction >> /dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/zkt-attendance/storage/logs/scheduler.log
